# eb-app.target
[Unit]
Description=EB App

[Install]
WantedBy=multi-user.target

